# security/opendnssec/pkg/ods_signerd.rc
#!/bin/sh
#
# $OpenBSD$

daemon="${TRUEPREFIX}/sbin/ods-signerd"

. /etc/rc.d/rc.subr

rc_cmd $1
